# decoys/smb/supervisord.conf
; The SMB decoy needs three processes cooperating:
;
;   rsyslogd     receives smbd_audit lines on local7, writes samba-audit.log
;   smbd         the actual SMB protocol implementation
;   opencanaryd  tails samba-audit.log and emits JSON events to stdout
;
; Ordering matters: rsyslogd must own the audit file before smbd emits to it.
; priority controls start order (lower starts first).

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root

[program:rsyslogd]
command=/usr/sbin/rsyslogd -n -f /etc/rsyslog.conf
priority=10
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:smbd]
command=/usr/sbin/smbd --foreground --no-process-group --debug-stdout
priority=20
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:opencanaryd]
command=/opt/venv/bin/opencanaryd --dev
priority=30
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
